# deploy/motorsports-digest.service
[Unit]
Description=F1/IndyCar morning buzz digest
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
Group=__USER__
WorkingDirectory=__DIR__
EnvironmentFile=__DIR__/.env
ExecStart=__DIR__/.venv/bin/python -m digest.main --config __DIR__/config.toml
